# deploy/systemd/unsubscribeme.service
[Unit]
Description=UnsubscribeMe Telegram YouTube RSS Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Пользователь, под которым работает бот (рекомендуется отдельный)
User=unsubscribeme
Group=unsubscribeme

# Рабочая директория — путь до клон‑репозитория на сервере
WorkingDirectory=/home/unsubscribeme/unsubscribeme

# Файл окружения с токеном и настройками
EnvironmentFile=/home/unsubscribeme/unsubscribeme/.env

# Путь до исполняемого файла в venv
ExecStart=/home/unsubscribeme/unsubscribeme/.venv/bin/unsubscribeme
Restart=always
RestartSec=5s

# Логи идут в journald: смотрите `journalctl -u unsubscribeme -f`
StandardOutput=journal
StandardError=journal

# Немного увеличим лимиты (по желанию)
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
